# rc.local
sudo sh PATH_TO_SCRIPT_VUE/script-vue.sh &
